# fedora-update/fedora-update-check.service
[Unit]
Description=Check for Fedora updates and notify

[Service]
Type=oneshot
ExecStart=/usr/bin/fedora-update-check.sh

[Install]
WantedBy=default.target
